# romeo-slack.service
[Unit]
Description=Romeo – The Boondock Doorbell

[Service]
EnvironmentFile=/etc/romeo
ExecStart=/usr/bin/romeo-slack.sh
