# romeo.service
[Unit]
Description=Romeo – The Boondock Dorbell

[Service]
EnvironmentFile=/etc/romeo
ExecStart=/usr/bin/romeo.sh
